# bin/deployments/systemd/atm0s-sdn-node.service
[Unit]
Description=atm0s-sdn-node service
After=network.target

[Service]
ExecStart=/opt/atm0s-sdn-node.sh
StandardOutput=append:/var/log/atm0s-sdn-node.log
StandardError=append:/var/log/atm0s-sdn-node.log

[Install]
WantedBy=multi-user.target